# server/src/services/honeycomb-verification.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
	Copyright 2007 Sun Microsystems, Inc.  All rights reserved.
	Use is subject to license terms.

	$Id: honeycomb-verification.xml 10855 2007-05-19 02:54:08Z bberndt $

	Service manifest for the Honeycomb Verification Service.
-->

<service_bundle type='manifest' name='SUNWhcserver:honeycomb-verification'>

<service
	name='application/honeycomb-verification'
	type='service'
	version='1'>

        <create_default_instance enabled='true' />

        <single_instance />

	<dependency
		name='honeycomb-configuration'
		grouping='require_all'
		restart_on='none'
		type='service'>
		<service_fmri value='svc:/application/honeycomb-configuration' />
	</dependency>

	<exec_method
		type='method'
		name='start'
		exec='/opt/honeycomb/bin/node_verification -a -d'
		timeout_seconds='-1' />

        <exec_method
                type='method'
                name='stop'
                exec=':true'
                timeout_seconds='3' />

        <property_group name='startd' type='framework'>
                <propval name='duration' type='astring' value='transient' />
        </property_group>

	<stability value='Unstable' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
				Honyecomb Verification
			</loctext>
		</common_name>
		<description>
			<loctext xml:lang='C'>
This is the Honeycomb Verification Service.
			</loctext>
		</description>
		<documentation>
			<manpage title='honeycomb-verification' section='1M'
				manpath='/usr/share/man' />
		</documentation>
	</template>
</service>

</service_bundle>
